# remont-web.service
[Unit]
Description=Remont Web Application (FastAPI)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/Remont
Environment="PATH=/path/to/Remont/venv/bin"
ExecStart=/path/to/Remont/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

